# BOX4s-main/crontab/amadmin.crontab
MAILTO=box@4sconsult.de
# m h  dom mon dow   command

#nmapd short scan
0 * * * *  $BASEDIR/$GITDIR/../listscan.sh | $BASEDIR$GITDIR/Scripts/System_Scripts/sendmail.sh box@4sconsult.de "Nmap List-Scan Log"

# Check Disk Space Daily
@daily  $BASEDIR/$GITDIR/scripts/Automation/diskspacealert.sh | $BASEDIR$GITDIR/Scripts/System_Scripts/sendmail.sh box@4sconsult.de "Festplattenspeicher erschöpft"

# Execute Curator Daily
@daily curator --config ~/.curator.yml ~/.actions.yml | $BASEDIR$GITDIR/Scripts/System_Scripts/sendmail.sh box@4sconsult.de "Elastic Curator Log"

# Score calculation
*/5 * * * * $BASEDIR/$GITDIR/Scripts/Automation/score_calculation/calculate_scores.sh
