# BOX4s-main/crontab/amadmin.crontab
MAILTO=box@4sconsult.de
# m h  dom mon dow   command

# Check Disk Space Daily
@daily  $BASEDIR/$GITDIR/Scripts/Automation/diskspacealert.sh | $BASEDIR$GITDIR/Scripts/System_Scripts/sendmail.sh box@4sconsult.de "Festplattenspeicher erschöpft"

# Execute Curator Daily
@daily $HOME/.local/bin/curator --config ~/curator.yml ~/actions.yml | $BASEDIR$GITDIR/Scripts/System_Scripts/sendmail.sh box@4sconsult.de "Elastic Curator Log"
