# config/crontab/amadmin.crontab
MAILTO=box@4sconsult.de
# m h  dom mon dow   command

# Check Disk Space Daily
@daily  $BASEDIR/$GITDIR/scripts/Automation/diskspacealert.sh
